# gumstix-oe/com.gumstix.collection/packages/udev/udev-118/85-sdmod.rules
# do not edit this file, it will be overwritten on update

ACTION!="add", GOTO="sdmod_end"

SUBSYSTEM=="scsi_device", SYSFS{type}=="0|7|14", RUN+="/sbin/modprobe sd_mod"

LABEL="sdmod_end"

